# presentation-app/docker/services/laravel-reverb/run
#!/command/execlineb -P
with-contenv
cd /var/www/html
php artisan reverb:start --debug --port=6001
