# usb-autoinstall.service
[Unit]
Description=Auto-install presentation app from USB if present
After=local-fs.target
ConditionPathExists=/media/usb/presentation-usb/install.sh

[Service]
Type=oneshot
ExecStart=/bin/bash /media/usb/presentation-usb/install.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
