# docker/services/web/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid www-data
php-fpm &
nginx -g "daemon off;"
